# config/crontab
# Update the SSL certificates once a month
0 0 1 * * /bin/bash /var/www/config/certbot.sh
